# DNCB440.d/backup/etc/cron.hourly/inn-cron-rnews
#!/bin/sh
/sbin/chkconfig innd || exit 0
su - news -c 'unset LANG; unset LC_COLLATE; /usr/lib/news/bin/rnews -U'
