# DNCB440.d/backup/etc/cron.weekly/mailsoda.cron
#!/bin/sh
mysqldump sodalord | mail -s "weekly sodalord mysql dump" sodalord-internal@mit.edu
